# datasets/splunk-detections/linux-adding-crontab-using-list-parameter/dataset.yaml
metadata:
  name: linux-adding-crontab-using-list-parameter
  version: 9.0.0
  description: The following analytic detects suspicious modifications to cron jobs
    on Linux systems using the crontab command with list parameters. It leverages
    data from Endpoint Detection and Response (EDR) agents, focusing on process names
    and command-line executions. This activity is significant as it may indicate an
    attempt to establish persistence or execute malicious code on a schedule. If confirmed
    malicious, the impact could include unauthorized code execution, data destruction,
    or other damaging out
  author: Teoderick Contreras, Bhavin Patel, Splunk
  tags:
  - hunting
  - splunk-security-content
  - cisco-isovalent-suspicious-activity
  - industroyer2
  - linux-privilege-escalation
  mitre_attack:
    techniques:
    - id: T1053.003
      name: Technique T1053.003
    tactics: []
files:
  bundled: []
  references:
  - uri: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/attack_techniques/T1053.003/crontab_list_parameter/sysmon_linux.log
    description: Syslog:Linux-Sysmon/Operational - sysmon:linux
    format: auto
    schema: raw
    sourcetype: sysmon:linux
  - uri: https://media.githubusercontent.com/media/splunk/attack_data/master/datasets/cisco_isovalent/cisco_isovalent.log
    description: not_applicable - cisco:isovalent:processExec
    format: auto
    schema: raw
    sourcetype: cisco:isovalent:processExec
defaults:
  replay:
    delta_factor: 1.0
    base_time: earliest
    target_time: now-1h
  schema: raw
